# services/couchdb/run
#!/bin/sh
exec /docker-entrypoint.sh /opt/couchdb/bin/couchdb
