# services/couchdb/run
#!/bin/sh
exec tini -- /docker-entrypoint.sh /opt/couchdb/bin/couchdb