# services/memcached/run
#!/bin/sh
exec docker-entrypoint.sh memcached -u memcache -m 256
