# services/dynamodb/run
#!/bin/sh
set -e
cd /home/dynamodblocal
exec java -jar DynamoDBLocal.jar -sharedDb -dbPath ./data
