# services/percona/run
#!/bin/sh
exec /docker-entrypoint.sh mysqld --user=root