# services/solr/finish
#!/bin/sh

/opt/solr/bin/solr stop -p ${SOLR_PORT:-8983}
