# services/solr/run
#!/bin/sh
exec docker-entrypoint.sh solr-foreground -force
