# services/solr/run
#!/bin/sh
chpst -u solr docker-entrypoint.sh solr-foreground
