# services/valkey/run
#!/bin/sh
exec docker-entrypoint.sh valkey-server
