# services/phpmyadmin/run
#!/bin/sh
cd /var/www/html
exec /run.sh supervisord -n -j /supervisord.pid
