# services/phpmyadmin/run
#!/bin/sh
exec /run.sh supervisord -n -j /supervisord.pid
