# services/mariadb/run
#!/bin/sh
exec docker-entrypoint.sh mariadbd
